# usr/lib/systemd/user/mcgeecloud.mount
[Unit]
Description=Mount for mcgeecloud NFS share for %u
After=network.target

[Mount]
What=nas:/var/mnt/bender/cloud/%u
Where=~/.mnt/mcgeecloud
Type=nfs4
Options=sec=krb5p

[Install]
WantedBy=default.target
