# usr/lib/systemd/user/mcgeecloud.automount
[Unit]
Description=User Automount for mcgeecloud NFS share

[Automount]
Where=%h/.mnt/mcgeecloud
TimeoutIdleSec=60

[Install]
WantedBy=default.target
