# usr/lib/systemd/user/mcgeecloud.automount
[Unit]
Description=Automount NFS Share for %i's mcgeecloud
Requires=network-online.target
After=network-online.target

[Automount]
Where=/home/%i/.mnt/mcgeecloud

[Install]
WantedBy=default.target
