# usr/lib/systemd/user/mcgeecloud.automount
[Unit]
Description=User Automount for mcgeecloud NFS share

[Automount]
Where=/var/mnt/mcgeecloud/%u
TimeoutIdleSec=60

[Install]
WantedBy=default.target
