# usr/lib/systemd/user/mcgeecloud.automount
[Unit]
Description=Automount NFS Share for mcgeecloud
Requires=network-online.target
After=network-online.target

[Automount]
Where=%h/.mnt/mcgeecloud

[Install]
WantedBy=default.target
