# usr/lib/systemd/user/var-mnt-mcgeecloud@\-.automount
[Unit]
Description=Mount NFS share for user %i
After=network.target

[Mount]
What=nas:/var/mnt/bender/shares/%i
Where=/var/mnt/mcgeecloud/shares/%i
Type=nfs4
Options=rw,timeo=100,sec=krb5p,fsc,x-gvfs-show

[Install]
WantedBy=multi-user.target