# usr/lib/systemd/system/home-%u-mcgeecloud.automount
[Unit]
Description=Automount NFS for mcgeecloud in /home/%u/mcgeecloud
Requires=network.target
After=network-online.target

[Automount]
Where=/home/%u/mcgeecloud

[Install]
WantedBy=multi-user.target